# files/clustersetup.path
[Path]
PathExists=/var/home/core/pull-secret.txt

[Install]
WantedBy=multi-user.target

